# services/sdwan-api.service
[Unit]
Description=SD-WAN API Server
After=network-online.target wg0.service
Requires=wg0.service

[Service]
Type=simple
ExecStart=/etc/sdwan/api/server.sh
Restart=always
RestartSec=10
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sdwan-api

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/sdwan

[Install]
WantedBy=multi-user.target
